# systemd/var-run-hadoop.path
[Unit]
Description=Hadoop rundir

[Path]
PathExists=/var/run/hadoop
